# overlay/houdini/tricks/CpuShares/supervisord.conf
[supervisord]
nodaemon=true

[program:program1]
command=/cpu_shares.py
autostart=true
autorestart=true

[program:program2]
command=/path/to/program2
autostart=true
autorestart=true
